# app/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:myservice]
command=python /usr/src/app/service.py -conf /usr/src/app/conf.yml